# systemd/system/pipe-bot.socket
[Socket]
ListenFIFO=%t/discord
Service=pipe-bot.service
